# data/autolearn-daily.service
[Unit]
Description=Claw Autonomous Learning Service
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/config/clawd
ExecStart=/usr/local/bin/autolearn
StandardOutput=append:/config/clawd/logs/autolearn_cron.log
StandardError=append:/config/clawd/logs/autolearn_cron.log

[Install]
WantedBy=multi-user.target
